# deploy/supervisor/qbazaar-reverb.conf
; /etc/supervisor/conf.d/qbazaar-reverb.conf
;
; Laravel Reverb (WebSocket server). Wired into the messaging flow in Sprint 8.

[program:qbazaar-reverb]
process_name=%(program_name)s
command=php /home/sanad/qbazaar/qbazaar-api/artisan reverb:start --host=127.0.0.1 --port=8080
autostart=true
autorestart=true
user=sanad
redirect_stderr=true
stdout_logfile=/var/log/supervisor/qbazaar-reverb.log
stopwaitsecs=10
