# deploy/supervisor/qbazaar-queue.conf
; /etc/supervisor/conf.d/qbazaar-queue.conf
;
; Laravel queue worker for QBazaar.
; On Linux we get pcntl + posix so this is the canonical worker.
; (On Windows we fall back to `php artisan queue:work` in a terminal — see qbazaar-api/DEV-SETUP.md.)

[program:qbazaar-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /home/sanad/qbazaar/qbazaar-api/artisan queue:work --queue=critical,default,media,search,low --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=sanad
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/qbazaar-queue.log
stopwaitsecs=3600
